# deploy/semantic-search-ingest.service
[Unit]
Description=Semantic search daily incremental ingestion
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/semantic-search
Environment=PROJECT_DIR=/opt/semantic-search
Environment=COMPANY_ID=gainr
ExecStart=/opt/semantic-search/scripts/run_scheduled_ingestion.sh
# A genuine full re-embedding can exceed one day on the 8 GB CPU host. The
# script lock prevents overlap, and its signal trap removes the run container.
TimeoutStartSec=48h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
